# packaging/linux/debian/jtnt-agentd.service
[Unit]
Description=JTNT RMM Agent
Documentation=https://docs.jtnt.us/agent
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/jtnt-agent

[Service]
Type=simple
ExecStart=/usr/local/bin/jtnt-agentd
Restart=always
RestartSec=10
User=jtnt-agent
Group=jtnt-agent

# Working directory
WorkingDirectory=/var/lib/jtnt-agent

# Environment
Environment="JTNT_STATE_DIR=/var/lib/jtnt-agent"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jtnt-agentd

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/jtnt-agent
ReadOnlyPaths=/etc/jtnt-agent

# Capability restrictions (none needed for standard operation)
CapabilityBoundingSet=
AmbientCapabilities=

# If network binding is needed, uncomment:
# CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# AmbientCapabilities=CAP_NET_BIND_SERVICE

# Namespace isolation
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
TasksMax=512

# Memory protection
MemoryDenyWriteExecute=true

# Filesystem protection
ProtectKernelLogs=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid

[Install]
WantedBy=multi-user.target
