# packaging/linux/debian/prerm
#!/bin/bash
# JTNT Agent - Pre-Remove Script for Debian/Ubuntu
# Runs before package removal

set -e

# Stop and disable the service
if [ -d /run/systemd/system ]; then
    if systemctl is-active --quiet jtnt-agentd.service; then
        echo "Stopping JTNT Agent service..."
        systemctl stop jtnt-agentd.service >/dev/null 2>&1 || true
    fi
    
    if systemctl is-enabled --quiet jtnt-agentd.service; then
        echo "Disabling JTNT Agent service..."
        systemctl disable jtnt-agentd.service >/dev/null 2>&1 || true
    fi
fi

exit 0
